# ldm.service.in
[Unit]
Description=lightweight device mounter

[Service]
EnvironmentFile=/etc/ldm.conf
ExecStart=@@BINDIR@@/ldm -u ${MOUNT_OWNER} -p ${BASE_MOUNTPOINT} -m ${FMASK_DMASK}
KillMode=process

[Install]
WantedBy=multi-user.target
