# ldm.service
[Unit]
Description=lightweight device mounter

[Service]
EnvironmentFile=/etc/conf.d/ldm
ExecStart=/usr/bin/ldm -u $USER_UID -g $USER_GID

[Install]
WantedBy=multi-user.target
