# ldm.service
[Unit]
Description=lightweight device mounter

[Service]
EnvironmentFile=/etc/ldm.conf
ExecStart=/usr/bin/ldm -u ${MOUNT_OWNER} -p ${BASE_MOUNTPOINT}
KillMode=process

[Install]
WantedBy=multi-user.target
